# scripts/systemd/offerclaw-web.service
[Unit]
Description=OfferClaw Next.js (port 3000)
After=offerclaw-api.service

[Service]
Type=simple
WorkingDirectory=/home/inuyasha/.openclaw/workspace/offerclaw/frontend
ExecStart=/usr/bin/npm run dev
Restart=always
RestartSec=8
StandardOutput=append:/tmp/openclaw-demo/web.log
StandardError=append:/tmp/openclaw-demo/web.log
Environment="HOME=/home/inuyasha"
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/home/inuyasha/.npm-global/bin"
Environment="PORT=3000"

[Install]
WantedBy=default.target
