# scripts/systemd/openclaw-dashboard.service
[Unit]
Description=OpenClaw Demo Dashboard (port 6688)
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/inuyasha/.openclaw/workspace/demo/dashboard.py
Restart=always
RestartSec=5
StandardOutput=append:/tmp/openclaw-demo/dashboard.log
StandardError=append:/tmp/openclaw-demo/dashboard.log
Environment="HOME=/home/inuyasha"

[Install]
WantedBy=default.target
